# .crontab
0 9 * * * /bin/bash -l -c 'cd ~/Sites/the_rotten_pirate/current && bundle exec rake execute --silent >> /dev/null 2>&1'
0,5,10,15,20,25,30,35,40,45,50,55 * * * * /bin/bash -l -c 'cd /Users/james/Sites/the_rotten_pirate/current && RAILS_ENV=production bundle exec rake download_from_watch_file --silent >> /dev/null 2>&1'
0,5,10,15,20,25,30,35,40,45,50,55 * * * * /bin/bash -l -c 'cd /Users/james/Sites/craigs_warn && ruby warner.rb >> tmp/log.txt 2>&1'
0 0 1,8,15,22 * * /bin/bash -l -c '/Users/jhart/Sites/rsync_backup/mac_mini-backup.sh'
#0,30 * * * * /bin/bash -l -c 'cd /Users/james/Sites/craigs_warn && ruby apt_warner.rb >> tmp/log.txt 2>&1'
*/5 * * * * /bin/bash -l -c 'cd ~/Sites/urban_spoon/current && ruby restaurant_notifier.rb >> tmp/cron.log 2>&1'